# DVR/Service/shutdown.service
[Unit]
Description=shutdown
After=multi-user.target

[Service]
Type=idle
Restart=always
RestartSec=3
ExecStart=/usr/bin/sudo /usr/bin/python3 /home/pi/DVR/Scripts/autoshutdown.py

[Install]
WantedBy=multi-user.target